# shared/roles/unison/files/com.workstation.wsync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.workstation.wsync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>$HOME/.local/bin/wsync</string>
    </array>
    <key>StartInterval</key>
    <!-- 300 seconds = 5 minutes -->
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/wsync-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/wsync-stderr.log</string>
</dict>
</plist>
