# shared/roles/backups/templates/backup-staleness.service.j2
[Unit]
Description=Check restic backup staleness

[Service]
Type=oneshot
ExecStart={{ ansible_facts.env.HOME }}/.local/bin/backup-staleness-check
