# shared/roles/backups/templates/backrest.service.j2
[Unit]
Description=Backrest restic backup manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/backrest
Environment=BACKREST_PORT={{ backrest_port }}
Environment=BACKREST_CONFIG={{ ansible_facts.env.HOME }}/.config/backrest/config.json
Environment=BACKREST_DATA={{ ansible_facts.env.HOME }}/.local/share/backrest
Environment=BACKREST_RESTIC_COMMAND={{ ansible_facts.env.HOME }}/.local/bin/restic-wrapper
Environment=XDG_CACHE_HOME={{ ansible_facts.env.HOME }}/.cache
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
