# shared/roles/backups/templates/com.user.backrest-wake.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.backrest-wake</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{ ansible_facts.env.HOME }}/.local/bin/backrest-wake-restart</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>{{ ansible_facts.env.HOME }}/.local/log/backrest-wake-restart.log</string>
    <key>StandardErrorPath</key>
    <string>{{ ansible_facts.env.HOME }}/.local/log/backrest-wake-restart.log</string>
</dict>
</plist>
