# shared/roles/sync-monitor/templates/sync-staleness.service.j2
[Unit]
Description=Check sync hub health

[Service]
Type=oneshot
ExecStart={{ ansible_facts.env.HOME }}/.local/bin/sync-health-check
